# deploy/supervisor-genomics.conf
[program:genomics-app]
; UPDATE THESE PATHS to your installation
command=python3 -m uvicorn backend.main:app --host 0.0.0.0 --port 8600
directory=/home/YOUR_USER/genomics
user=YOUR_USER
autostart=true
autorestart=true
stderr_logfile=/var/log/genomics-app.err.log
stdout_logfile=/var/log/genomics-app.out.log
environment=HOME="/home/YOUR_USER",PYTHONPATH="/home/YOUR_USER/genomics"
